# providers/os/resources/services/testdata/systemd/etc/systemd/system/fstrim.timer
[Unit]
Description=Discard unused blocks once a week

[Timer]
OnCalendar=weekly
Persistent=yes
Unit=fstrim.service

[Install]
WantedBy=timers.target
