# providers/os/resources/services/testdata/systemd/etc/systemd/system/logrotate.timer
[Unit]
Description=Daily rotation of log files

[Timer]
OnCalendar=daily
Persistent=yes

[Install]
WantedBy=timers.target
