# providers/os/resources/services/testdata/systemd/etc/systemd/system/static-check.timer
[Unit]
Description=Static timer without Install section

[Timer]
OnCalendar=hourly
